# .docker/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:laravel-scheduler]
command=sh -c "while true; do php /var/www/artisan schedule:run --verbose --no-interaction & sleep 60; done"
autostart=true
autorestart=true
user=www
stdout_logfile=/var/log/supervisor/scheduler.log
stderr_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=2

; ============================================
; Debezium Outbox Consumer (один процесс)
; ============================================
[program:debezium-consumer]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan kafka:consume-outbox
autostart=true
autorestart=true
numprocs=1
user=www
stdout_logfile=/var/log/supervisor/debezium-consumer-%(process_num)s.log
stderr_logfile=/var/log/supervisor/debezium-consumer-%(process_num)s.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=3

; ============================================
; Outbox Fallback (если Debezium упал)
; ============================================
[program:outbox-fallback]
command=sh -c "while true; do php /var/www/artisan outbox:process-pending; sleep 10; done"
autostart=true
autorestart=true
user=www
stdout_logfile=/var/log/supervisor/outbox-fallback.log
stderr_logfile=/var/log/supervisor/outbox-fallback.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=4

; ============================================
; Kafka Consumers для нотификаций
; ============================================

; Приоритет: transactional (критичные) - больше всего воркеров
[program:kafka-notifications-transactional]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan kafka:consume-notifications transactional
autostart=true
autorestart=true
numprocs=4
user=www
stdout_logfile=/var/log/supervisor/kafka-transactional-%(process_num)s.log
stderr_logfile=/var/log/supervisor/kafka-transactional-%(process_num)s.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=5

; Приоритет: high
[program:kafka-notifications-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan kafka:consume-notifications high
autostart=true
autorestart=true
numprocs=3
user=www
stdout_logfile=/var/log/supervisor/kafka-high-%(process_num)s.log
stderr_logfile=/var/log/supervisor/kafka-high-%(process_num)s.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=6

; Приоритет: normal - больше всего сообщений
[program:kafka-notifications-normal]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan kafka:consume-notifications normal
autostart=true
autorestart=true
numprocs=8
user=www
stdout_logfile=/var/log/supervisor/kafka-normal-%(process_num)s.log
stderr_logfile=/var/log/supervisor/kafka-normal-%(process_num)s.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=7

; Приоритет: low (маркетинговые) - меньше всего воркеров
[program:kafka-notifications-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan kafka:consume-notifications low
autostart=true
autorestart=true
numprocs=2
user=www
stdout_logfile=/var/log/supervisor/kafka-low-%(process_num)s.log
stderr_logfile=/var/log/supervisor/kafka-low-%(process_num)s.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=8
